# crontab
*/5 * * * * root sh /scripts/run_certbot.sh
